# etc/serve_https.socket
[Unit]
Description=serve HTTPS socket

[Socket]
ListenStream=443
Service=serve
