# etc/serve_http.socket
[Unit]
Description=serve HTTP socket

[Socket]
ListenStream=80
Service=serve
